# profile.d/irbrc.sh
[ -f ~/.irbrc ] || cp /etc/skel/.irbrc ~/.irbrc
