# profile.d/pythonrc.sh
[ -e ~/.pythonrc ] || cp /etc/skel/.pythonrc ~/.pythonrc
export PYTHONSTARTUP=~/.pythonrc
