# profile.d/vimrc.sh
[ -f ~/.vimrc ] || cp /etc/skel/.vimrc ~/.vimrc
