# etc/profile.d/oh-my-zsh.sh
if [ $USER != 'root' ]; then
  [ -d ~/.oh-my-zsh ] || /opt/oh-my-zsh/tools/install.sh
fi
