# etc/profile.d/pythonrc.sh
[ -e ~/.pythonrc ] && export PYTHONSTARTUP=~/.pythonrc
